# runit/.config/sv/udiskie/run
#!/bin/sh
exec udiskie --automount --notify --tray
