# runit/.config/sv/wlsunset/run
#!/bin/sh

exec wlsunset -l 39.7 -L "-105" -t 3500
