# runit/.config/sv/swaybg/run
#!/bin/sh
exec swaybg -i $HOME/.local/share/wallpapers/current-wallpaper -m fill

